# compose/docker/proxy/s6-rc.d/init-certs/up
#!/command/execlineb -P
/usr/local/bin/generate-certs.sh

