# compose/docker/proxy/s6-rc.d/traefik/run
#!/command/with-contenv bash
exec /usr/local/bin/traefik --configFile=/etc/traefik/traefik.yml

